# mbuswebd/mbuswebd.service
[Unit]
Description=M-Bus  WEB-server daemon
After=network.target
Before=nginx.service

[Service]
ExecStart=PYTHONNAME SERVICENAME
Type=simple
Restart=always


[Install]
WantedBy=default.target
RequiredBy=network.target
